# examples/systemd/proton-workflow-connector.service
[Unit]
Description=Proton Workflow Connector
After=network-online.target protonmail-bridge-headless.service
Wants=network-online.target protonmail-bridge-headless.service

[Service]
Type=simple
WorkingDirectory=%h/proton-workflow-connector
ExecStart=%h/proton-workflow-connector/.venv/bin/proton-workflow-connector --transport streamable-http --host 127.0.0.1 --port 8765 --env-file %h/.config/proton-workflow-connector/env
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
NoNewPrivileges=true

[Install]
WantedBy=default.target
